# config/home/user/.config/systemd/user/profile-cleaner.timer
[Unit]
Description=Run profile-cleaner daily

[Timer]
OnCalendar=daily
Persistent=true
RandomizedDelaySec=2

[Install]
WantedBy=timers.target
