# private_dot_config/autostart/protonmail-bridge.desktop
[Desktop Entry]
Type=Application
Name=ProtonMail Bridge
Exec="/usr/bin/protonmail-bridge" "--no-window"
X-GNOME-Autostart-enabled=true
